# route53-dns-updater.service.template
[Unit]
Description=Route 53 DNS Updater

[Service]
Type=simple
ExecStart=/path/to/script
Description=Route 53 DNS Updater
